# scripts/nido-backup.timer
[Unit]
Description=Daily Nido SQLite backup at 3:00 UTC

[Timer]
OnCalendar=*-*-* 03:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
